# uu.conf
[program:uu]
command=python3 /home/curator2/uu/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/uu/test.err.log
stdout_logfile=/home/curator2/uu/test.out.log